# deploy/beeroot/wow-dashboard-deploy-poll.timer
[Unit]
Description=Check for WoW Dashboard production deployment requests

[Timer]
OnBootSec=30s
OnUnitActiveSec=30s
AccuracySec=5s
Unit=wow-dashboard-deploy-poll.service

[Install]
WantedBy=timers.target
